# deploy/systemd/calb-ngrok@.service
[Unit]
Description=ngrok tunnel for CALB SizingTool (%i)
After=network-online.target calb-sizingtool@%i.service
Wants=network-online.target calb-sizingtool@%i.service

[Service]
Type=simple
User=calb
EnvironmentFile=/etc/calb-sizingtool/%i.env
WorkingDirectory=/opt/calb/%i/CALB_SIZINGTOOL
Environment=NGROK_CONFIG=/home/calb/.config/ngrok/ngrok.yml
ExecStart=/usr/local/bin/calb-ngrok-wrapper %i
# Use conservative restart policy to avoid tight restart loops when ngrok errors occur
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=60
StartLimitBurst=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
